# config/systemd/user/gitwatch.service
[Unit]
Description=gitwatch

[Service]
ExecStart=/dotfiles/gitwatch/gitwatch.sh $HOME/workspace
Restart=always

[Install]
WantedBy=sway.target
