# config/systemd/user/gitwatch.service
[Unit]
Description=gitwatch
PartOf=graphical-session.target

[Service]
ExecStart=$HOME/dotfiles/gitwatch/gitwatch.sh $HOME/workspace
Restart=always

[Install]
WantedBy=sway.target
