# config/systemd/user/gitwatch.service
[Unit]
Description=gitwatch

[Service]
ExecStart=$HOME/dotfiles/gitwatch/gitwatch.sh $HOME/workspace
Restart=always

[Install]
WantedBy=sway.target
