# config/systemd/user/gitwatch.service
[Unit]
Description=xmod
PartOf=graphical-session.target

[Service]
ExecStart=
# Restart=always

[Install]
WantedBy=xsession.target
