# config/systemd/user/gitwatch.service
[Unit]
Description=gitwatch

[Service]
ExecStart=/home/yasha/dotfiles/gitwatch/gitwatch.sh /home/yasha/workspace
Restart=always

[Install]
WantedBy=default.target
