# config/systemd/user/xmod.service
[Unit]
Description=xmod
PartOf=graphical-session.target

[Service]
ExecStart=/usr/sbin/xmodmap /home/yasha/.Xmodmap
Restart=never

[Install]
WantedBy=xsession.target
